# ddev-drupal-core/scripts/drupal-cache-updater.service
[Unit]
Description=Update drupal-core seed cache for Coder workspaces
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
# Must run as the user who owns the seed cache directory — NOT root.
# DDEV refuses to run as root; this service must run as a normal user.
# Change User and WorkingDirectory to match the account that owns the cache.
User=rfay
WorkingDirectory=/home/rfay
# The script is installed to /usr/local/bin/ (see docs/admin/server-setup.md).
# Pass --seed-dir if your seed directory differs from the default.
ExecStart=/usr/local/bin/update-drupal-cache --seed-dir /home/rfay/cache/drupal-core-seed
StandardOutput=journal
StandardError=journal
# Allow up to 30 minutes for composer update + site install
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
